# var/lib/dpkg/info/incron.list
/.
/etc
/etc/incron.allow
/etc/incron.conf
/etc/incron.d
/etc/incron.deny
/etc/init.d
/etc/init.d/incron
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/incron.service
/usr
/usr/bin
/usr/bin/incrontab
/usr/sbin
/usr/sbin/incrond
/usr/share
/usr/share/doc
/usr/share/doc/incron
/usr/share/doc/incron/README
/usr/share/doc/incron/README.Debian
/usr/share/doc/incron/TODO
/usr/share/doc/incron/changelog.Debian.amd64.gz
/usr/share/doc/incron/changelog.Debian.gz
/usr/share/doc/incron/changelog.gz
/usr/share/doc/incron/copyright
/usr/share/doc/incron/examples
/usr/share/doc/incron/examples/incron.conf.example
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/incrontab.1.gz
/usr/share/man/man5
/usr/share/man/man5/incron.conf.5.gz
/usr/share/man/man5/incrontab.5.gz
/usr/share/man/man8
/usr/share/man/man8/incrond.8.gz
/var
/var/spool
/var/spool/incron
