# var/lib/dpkg/info/ntp.list
/.
/etc
/etc/NetworkManager
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/ntp
/etc/apparmor
/etc/apparmor/init
/etc/apparmor/init/network-interface-security
/etc/apparmor.d
/etc/apparmor.d/tunables
/etc/apparmor.d/tunables/ntpd
/etc/apparmor.d/usr.sbin.ntpd
/etc/cron.daily
/etc/cron.daily/ntp
/etc/default
/etc/default/ntp
/etc/dhcp
/etc/dhcp/dhclient-exit-hooks.d
/etc/dhcp/dhclient-exit-hooks.d/ntp
/etc/init.d
/etc/init.d/ntp
/etc/ntp.conf
/usr
/usr/bin
/usr/bin/calc_tickadj
/usr/bin/ntpdc
/usr/bin/ntpq
/usr/bin/ntpsweep
/usr/bin/ntptrace
/usr/bin/sntp
/usr/bin/update-leap
/usr/sbin
/usr/sbin/ntp-keygen
/usr/sbin/ntp-wait
/usr/sbin/ntpd
/usr/sbin/ntptime
/usr/share
/usr/share/doc
/usr/share/doc/ntp
/usr/share/doc/ntp/NEWS.Debian.gz
/usr/share/doc/ntp/NEWS.gz
/usr/share/doc/ntp/README.Debian.gz
/usr/share/doc/ntp/README.refclocks
/usr/share/doc/ntp/README.versions
/usr/share/doc/ntp/changelog.Debian.gz
/usr/share/doc/ntp/copyright
/usr/share/doc/ntp/ntp-wait.html
/usr/share/doc/ntp/ntp.conf.html
/usr/share/doc/ntp/ntp.keys.html
/usr/share/doc/ntp/ntpd.html
/usr/share/doc/ntp/ntpdc.html
/usr/share/doc/ntp/ntpq.html
/usr/share/doc/ntp/ntpsnmpd.html
/usr/share/doc/ntp/ntpsweep.html
/usr/share/doc/ntp/ntptrace.html
/usr/share/doc/ntp/update-leap.html
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/calc_tickadj.1.gz
/usr/share/man/man1/ntp-wait.1.gz
/usr/share/man/man1/ntpdc.1.gz
/usr/share/man/man1/ntpq.1.gz
/usr/share/man/man1/ntpsweep.1.gz
/usr/share/man/man1/ntptrace.1.gz
/usr/share/man/man1/sntp.1.gz
/usr/share/man/man1/update-leap.1.gz
/usr/share/man/man5
/usr/share/man/man5/ntp.conf.5.gz
/usr/share/man/man5/ntp.keys.5.gz
/usr/share/man/man8
/usr/share/man/man8/ntp-keygen.8.gz
/usr/share/man/man8/ntp-wait.8.gz
/usr/share/man/man8/ntpd.8.gz
/usr/share/man/man8/ntptime.8.gz
/usr/share/ntp
/usr/share/ntp/lib
/usr/share/ntp/lib/NTP
/usr/share/ntp/lib/NTP/Util.pm
/var
/var/lib
/var/lib/ntp
/var/log
/var/log/ntpstats
/etc/apparmor/init/network-interface-security/usr.sbin.ntpd
