# var/lib/dpkg/info/pdns-server.list
/.
/etc
/etc/default
/etc/default/pdns
/etc/init.d
/etc/init.d/pdns
/etc/powerdns
/etc/powerdns/pdns.conf
/etc/powerdns/pdns.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pdns.service
/lib/systemd/system/pdns@.service
/usr
/usr/bin
/usr/bin/pdns_control
/usr/bin/pdnsutil
/usr/bin/zone2json
/usr/bin/zone2sql
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/pdns
/usr/lib/x86_64-linux-gnu/pdns/librandombackend.so
/usr/sbin
/usr/sbin/pdns_server
/usr/share
/usr/share/doc
/usr/share/doc/pdns-server
/usr/share/doc/pdns-server/NEWS.Debian.gz
/usr/share/doc/pdns-server/README.Debian
/usr/share/doc/pdns-server/changelog.Debian.gz
/usr/share/doc/pdns-server/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/pdns-server
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pdns_control.1.gz
/usr/share/man/man1/pdns_server.1.gz
/usr/share/man/man1/pdnsutil.1.gz
/usr/share/man/man1/zone2json.1.gz
/usr/share/man/man1/zone2sql.1.gz
