# var/lib/dpkg/info/pdns-server.conffiles
/etc/default/pdns
/etc/init.d/pdns
/etc/powerdns/pdns.conf
