# var/lib/dpkg/info/sshguard.conffiles
/etc/default/sshguard
/etc/init.d/sshguard
/etc/sshguard/whitelist
