# var/lib/dpkg/info/sshguard.list
/.
/etc
/etc/default
/etc/default/sshguard
/etc/init.d
/etc/init.d/sshguard
/etc/sshguard
/etc/sshguard/whitelist
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/sshguard.service
/usr
/usr/lib
/usr/lib/sshguard
/usr/lib/sshguard/firewall
/usr/lib/sshguard/sshguard-journalctl
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/sshg-fw
/usr/lib/x86_64-linux-gnu/sshg-logtail
/usr/lib/x86_64-linux-gnu/sshg-parser
/usr/sbin
/usr/sbin/sshguard
/usr/share
/usr/share/doc
/usr/share/doc/sshguard
/usr/share/doc/sshguard/changelog.Debian.gz
/usr/share/doc/sshguard/changelog.gz
/usr/share/doc/sshguard/copyright
/usr/share/doc/sshguard/examples
/usr/share/doc/sshguard/examples/whitelistfile.example
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/sshguard.8.gz
