# var/backups/inotify/etc/init.d/firewall_2018.12.12_01:52:29
#!/bin/bash

### BEGIN INIT INFO
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Description: firewall service
### END INIT INFO

# description: Firewall service for adding and removing rules

start() {
        echo "configuring firewall..."
        /bin/bash /etc/scripts/firewall_start
}

stop() {
        echo "stopping firewall..."
        /bin/bash /etc/scripts/firewall_stop
}

case "$1" in
        start)
                start
                ;;
        stop)
                stop
                ;;
        restart)
                stop
                start
                ;;
        *)
                start
esac
